# skills-setup/granola-sync/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.alpha.granola-sync</string>
  <key>ProgramArguments</key>
  <array>
    <!-- EDIT: replace USERNAME with your Mac username -->
    <string>/Users/USERNAME/.openclaw/workspace/scripts/granola-sync/watch.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/USERNAME/.openclaw/workspace/scripts/granola-sync/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/USERNAME/.openclaw/workspace/scripts/granola-sync/stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
